# deploy/ghenkins.service
[Unit]
Description=ghenkins CI daemon (containerized)
After=network-online.target podman.socket
Wants=network-online.target podman.socket

[Service]
Type=simple
# Run the container; the script handles socket detection and volume mounts.
ExecStart=%h/projects/ghenkins/deploy/run-container.sh
ExecStop=podman stop ghenkins
Restart=on-failure
RestartSec=5s
Environment=GHENKINS_LOG_LEVEL=info
Environment=GHENKINS_IMAGE=localhost/ghenkins:latest
# Allow up to 10 minutes for in-progress CI runs to finish on shutdown.
TimeoutStopSec=10min
KillMode=mixed

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
