# deploy/ghenkins.service
[Unit]
Description=ghenkins CI daemon
After=network-online.target
Wants=network-online.target
Wants=podman.socket
After=podman.socket

[Service]
Type=simple
ExecStart=%h/.local/bin/ghenkins serve
Restart=on-failure
RestartSec=5s
Environment=GHENKINS_LOG_LEVEL=info
# Allow up to 10 minutes for in-progress CI runs to finish on shutdown.
TimeoutStopSec=10min
# Send SIGTERM to the main process; let it finish gracefully before SIGKILL.
KillMode=mixed

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
